# lifecycle-report.service
[Unit]
Description=Report changes in product/package lifecycle

[Service]
Type=oneshot
ExecStart=/usr/share/lifecycle/lifecycle-report
